# systemd/keyfetcher.service
[Install]
WantedBy = multi-user.target

[Service]
ExecStart=/opt/keyserver/bin/keyfetcher
