# systemd/keyfetcher.service
[Install]
WantedBy = multi-user.target

[Service]
ExecStart=/usr/local/keyserver/bin/keyfetcher
